# files/mysqld_exporter.socket
[Unit]
Description=Mysqld Exporter

[Socket]
ListenStream=9104

[Install]
WantedBy=sockets.target